# deploy/hf-space/supervisord.conf
[supervisord]
nodaemon=true
logfile=%(ENV_TEXLYRE_DATA_DIR)s/logs/supervisor/supervisord.log
pidfile=/tmp/supervisor/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=redis-server --port 6380 --bind 127.0.0.1 --dir /tmp/redis --save "" --appendonly no --protected-mode yes
autostart=true
autorestart=true
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; PRE-FLIGHT REQUIRED:
; Verify this directory inside ruthenian8/texlyre-texlive-ondemand-server
; before accepting the final image. Expected command is `python3 wsgi.py`.
[program:texlive-server]
directory=/app
command=python3 wsgi.py
environment=PORT="3001",REDIS_PORT="6380",API_ORIGINS="*"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ywebrtc-server]
directory=/app/services/y-webrtc-server
command=node ./bin/server.js
environment=NODE_ENV="production",PORT="3002",API_ORIGINS="*"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:peerjs-server]
directory=/app/services/peerjs-server
command=node dist/bin/peerjs.js --port 3003 --path /texlyre --key peerjs --proxied true --allow_discovery false
environment=NODE_ENV="production",PORT="3003",PEERJS_PORT="3003",PEERJS_PATH="/texlyre",PEERJS_KEY="peerjs",PEERJS_PROXIED="true",PEERJS_ALLOW_DISCOVERY="false",PEERJS_CONCURRENT_LIMIT="5000"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:filepizza-server]
directory=/app/services/filepizza-server
command=node server.js
environment=NODE_ENV="production",PORT="3004",REDIS_URL="redis://127.0.0.1:6380",API_ORIGINS="*",PEERJS_SERVERS="/peerjs/texlyre"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
